# ops-aws-infrastructure/ec2-setup/bright-valley-api.service
[Unit]
Description=bright-valley API Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=bright-valley
Group=bright-valley
WorkingDirectory=/opt/bright-valley
ExecStart=/opt/bright-valley/ec2-setup/bright-valley-api.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bright-valley

# Environment variables
Environment=JAVA_HOME=/usr/lib/jvm/java-21-amazon-corretto
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=-/opt/bright-valley/env/spring.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/bright-valley

# Resource limits (optional - adjust as needed)
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target